# root/etc/services.d/openvpn/run
#!/usr/bin/with-contenv sh

openvpn \
  --config /ovpn/client.conf \
  --auth-nocache 
